# meeting-transcript-system/config/com.obsidian.synctranscripts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian.synctranscripts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/fallonjensen/Obsidian Vault/Automations/meeting-transcript-system/sync-transcripts.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/fallonjensen/Obsidian Vault/Automations/meeting-transcript-system/transcript-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/fallonjensen/Obsidian Vault/Automations/meeting-transcript-system/transcript-sync.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
