# claude-session-memory-system/config/com.user.claude-session-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-session-backup</string>

    <!-- IMPORTANT: Replace YOUR_USERNAME with your actual username -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/Obsidian Vault/Automations/claude-session-backup/backup-claude-sessions.sh</string>
    </array>

    <!-- Run every 900 seconds (15 minutes) -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- Run immediately on system load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- IMPORTANT: Replace YOUR_USERNAME with your actual username -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Obsidian Vault/Automations/claude-session-backup/launchagent-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Obsidian Vault/Automations/claude-session-backup/launchagent-stderr.log</string>

    <!-- Include Homebrew paths for jq and other tools -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
